# provisioning/ansible/roles/beanstalkd/files/beanstalkd.conf
[program:beanstalkd]
command=/usr/local/bin/beanstalkd -b /var/beanstalkd